# webservice/web/conf/web.conf
[program:web]
command="/bin/bash -c '. /code/web/.venv/bin/activate; python /code/web/app.py'"
